# godot-updater.service
[Unit]
Description=Godot Updater Service
After=network.target

[Service]
ExecStart=/path/to/godot_installer.sh
Restart=on-failure

[Install]
WantedBy=default.target
